# etc/crontab
@reboot /home/pi/garage-remote/etc/start.sh 2>/dev/null &
2 0 * * * /usr/bin/sudo /sbin/reboot
